# init/multi-user.target
# Multi-user target for booting into multi-user mode

[Dependencies]


[Service]
Bin=/usr/bin/login
Type=regular
