# usystem/devmgr/devmgr.service
# devmgr service

[Service]
Bin=/sbin/devmgrd
Type=regular
