# src/network/netctld/netctld.service
# netctld service

[Service]
Bin=/sbin/netctld
Type=regular
