# devmgr/devmgr.service
# devmgr service

[Service]
Bin=/sbin/devmgr
Type=regular
