# src/core/coredumpd/coredumpd.service
# Coredumping service

[Service]
Bin=/sbin/coredumpd
Type=regular
