# dhcpcd/dhcpcd.service
# dhcpcd service

[Service]
Bin=/sbin/dhcpcd
Type=regular
